# deploy/model-storage-manager/model-storage-manager.service
[Unit]
Description=MSM-06: model-storage-manager watchdog (reconcile + sweep + gc via Chord control API)
Documentation=https://git.example.com/moosenet/lumina-constellation
After=network.target chord.service
# Not "Requires=chord.service": if Chord is briefly down this run should still
# execute, fail its curl calls, alert, and exit non-zero — not be skipped by
# systemd. The timer's next tick retries.

[Service]
Type=oneshot
# No secrets here — the script reads CHORD_JWT_SECRET/CHORD_CONTROL_PORT from
# Chord's own env file (CHORD_ENV_FILE, default <path>/.env) AT RUNTIME.
Environment=CHORD_ENV_FILE=<path>/.env
Environment=CHORD_CONTROL_HOST=127.0.0.1
Environment=MSM_HEARTBEAT_FILE=/var/log/chord/model-storage-manager.jsonl
Environment=MSM_HIGH_WATER_PERCENT=90
ExecStart=<path>/deploy/model-storage-manager/model-storage-manager.sh
# This process only reads Chord's env file and talks to 127.0.0.1 — no need
# for elevated privileges. Run as whichever user chord.service itself runs as
# (verify with `systemctl show chord.service --property=User,Group` on the
# deploy host and set these to match) so it can read CHORD_ENV_FILE without a
# permissions workaround; placeholder below assumes that user is "chord".
User=chord
Group=chord
TimeoutStartSec=180
